# rsmonitor_0.0.1-1_all/etc/systemd/system/rsmonitor.service
[Unit]
Description=A simple express sample project

[Service]
ExecStart=/usr/bin/rsmonitor
Restart=always
Group=nogroup
WorkingDirectory=/etc/rsmonitor

[Install]
WantedBy=multi-user.target
